# templates/default/etc/profile.d/oracle.sh.erb
export ORACLE_HOME=<%= @home_path %>
export ORACLE_CLIENT_BIN_PATH=<%= @bin_path %>
export ORACLE_CLIENT_LIB_PATH=<%= @lib_path %>

export PATH=$PATH:$ORACLE_CLIENT_BIN_PATH
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$ORACLE_CLIENT_LIB_PATH
