# templates/default/etc/profile.d/oracle.sh.erb
export ORACLE_HOME=<%= @oracle_home_path %>
export ORACLE_CLIENT_BIN_PATH=<%= @oracle_client_bin_path %>
export ORACLE_CLIENT_LIB_PATH=<%= @oracle_client_lib_path %>

export PATH=$PATH:$ORACLE_CLIENT_BIN_PATH
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$ORACLE_CLIENT_LIB_PATH
